# servers/puppet/modules/zulip/files/cron.d/email-mirror
MAILTO=root

* * * * *   humbug cd /home/zulip/deployments/current && python manage.py email-mirror
