# servers/puppet/modules/zulip/files/cron.d/email-mirror
MAILTO=root

* * * * *   humbug cd /home/humbug/humbug-deployments/current && python manage.py email-mirror
