# servers/puppet/modules/zulip/files/cron.d/graphite_backup
MAILTO=root
SHELL=/bin/bash

0 3 * * * humbug /home/humbug/zulip/servers/puppet/modules/zulip/files/graphite/daily_rsync_backup.sh
