# servers/puppet/modules/zulip/files/supervisor/conf.d/feedback-bot.conf
[program:feedback-bot]
command=/home/zulip/zulip/bots/feedback-bot
autostart=true
user=humbug
