# zulip/puppet/zulip_internal/files/cron.d/log2zulip
MAILTO=root

* * * * *   zulip cd /home/zulip/deployments/current && ./bots/log2zulip
