# puppet/zulip_internal/files/cron.d/send-digest-emails
MAILTO=root

# Send digest emails once a day.
0 11 * * *   zulip cd /home/zulip/deployments/current && python manage.py enqueue_digest_emails
