# puppet/zulip/files/postgresql/postgresql@.service.d/zulip.conf
# This file is managed by Zulip's Puppet configuration.
#
# Linux 6.14+ charges (per-user) the rings from io_uring against to
# RLIMIT_MEMLOCK, which is by default limited to 8MiB; the default
# Zulip PostgreSQL configuration needs ~13MiB.  Set a limit sizably
# above that.
[Service]
LimitMEMLOCK=256M
