# puppet/zulip/files/cron.d/archive-messages
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

# Daily cron job for retention policy - archives existing messages based on the policy
# and cleans expired entries from the archive.
0 6 * * * zulip /home/zulip/deployments/current/manage.py archive_messages >/dev/null
# Daily cron job for cleaning up old, unclaimed attachments and their corresponding files.
# Ideally should not run at the same time as the previous job, to avoid races.
0 5 * * * zulip /home/zulip/deployments/current/manage.py delete_old_unclaimed_attachments -f >/dev/null
