# puppet/zulip/files/cron.d/archive-messages
MAILTO=root
# Removing old archived data should follow after archiving messages,
# cause it checks related archived messages and original attachments rows before deleting .

# Archive messages once a day. Time is in UTC.
0 1 * * *   zulip cd /home/zulip/deployments/current && ./manage.py archive_messages

# Remove old archived data. Time is in UTC.
0 2 * * *   zulip cd /home/zulip/deployments/current && ./manage.py remove_old_archived_data

