# puppet/zulip/files/cron.d/archive-messages
MAILTO=root
# The retention policy data deletion process should remove old archived
# data *after* archiving messages, to check related archived messages
# and the original attachments' rows before deleting.

# Archive messages once a day. Time is in UTC.
0 1 * * *   zulip cd /home/zulip/deployments/current && ./manage.py archive_messages

# Remove old archived data. Time is in UTC.
0 2 * * *   zulip cd /home/zulip/deployments/current && ./manage.py remove_old_archived_data

