# puppet/zulip/files/cron.d/certbot-renew
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=root

# Cron job to renew certbot twice a day.
52 0,12 * * * root /root/certbot-auto renew --webroot --webroot-path=/var/www/certbot/ --quiet
