# puppet/zulip/files/cron.d/downgrade-small-realms-behind-on-payments
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

0 17 * * * zulip /home/zulip/deployments/current/manage.py downgrade_small_realms_behind_on_payments
