# puppet/zulip/files/cron.d/fetch-tor-exit-nodes
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

17 * * * * zulip /home/zulip/deployments/current/manage.py fetch_tor_exit_nodes
