# puppet/zulip/files/cron.d/rabbitmq-numconsumers
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=root

* * * * * root /home/zulip/deployments/current/scripts/nagios/check-rabbitmq-consumers
