# puppet/zulip/files/cron.d/rabbitmq-numconsumers
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=root

* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file notify_tornado
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file user_activity
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file user_activity_interval
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file user_presence
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file invites
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file signups
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file message_sender
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file feedback_messages
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file error_reports
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file digest_emails
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file user_activity_test
* * * * * root /home/zulip/deployments/current/bin/write-rabbitmq-consumers-state-file missedmessage_mobile_notifications
