# puppet/zulip/files/cron.d/restart-zulip
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

# Cron job to restart the Zulip server weekly, with the goal of
# minimizing the impact of any memory leaks that we might grow.
0 6 * * 7 zulip /home/zulip/deployments/current/scripts/restart-server
