# puppet/zulip/files/cron.d/update-analytics-counts
MAILTO=root

# Run update analytics tables every 5 minutes past the hour
5 * * * *   zulip cd /home/zulip/analytics/management/commands && python manage.py update_analytics_counts

