# puppet/zulip_ops/files/cron.d/active-user-stats
MAILTO=root

*/10 * * * *   zulip cd /home/zulip/deployments/current && ./manage.py active_user_stats
