# puppet/zulip_ops/files/cron.d/clearsessions
MAILTO=root

# Clear all expired Django sessions at 10:22 PM every day.
22 22 * * *   zulip cd /home/zulip/deployments/current && ./manage.py clearsessions
