# puppet/zulip_ops/files/cron.d/fetch-contributor-data
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

0 8 * * * zulip /home/zulip/deployments/current/tools/fetch-contributor-data
