# puppet/zulip_ops/files/cron.d/graphite_backup
MAILTO=root
SHELL=/bin/bash

0 3 * * * zulip /home/zulip/zulip/puppet/zulip_ops/files/graphite/daily_rsync_backup.sh
