# puppet/zulip_ops/files/cron.d/log2zulip
MAILTO=root

* * * * *   zulip cd /home/zulip/deployments/current && ./zulip-current-venv/share/zulip/integrations/log2zulip/log2zulip
