# docker/zulip-puppet/files/supervisor/conf.d/zulip_postsetup.conf
[program:zulip-postsetup-create_user]
command = bash -c "sleep 2;/opt/createZulipAdmin.sh"
user=zulip
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = false
exitcodes = 100,200
