# zulip/puppet/zulip/files/cron.d/email-mirror
MAILTO=root

* * * * *   zulip cd /home/zulip/deployments/current && python manage.py email_mirror
