# zulip/puppet/zulip/files/cron.d/restart-zulip
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
USER=zulip

0 6 * * 7 zulip supervisorctl restart all
