# zulip/puppet/zulip/files/cron.d/send-digest-emails
MAILTO=root

# Send digest emails once a day. Time is in UTC.
0 18 * * *   zulip cd /home/zulip/deployments/current && python manage.py enqueue_digest_emails
