# puppet/zulip_ops/templates/supervisor/conf.d/grafana.conf.erb
[program:grafana]
command=<%= @bin %> -config /etc/grafana/grafana.ini
directory=<%= @dir %>
priority=15
autostart=true
autorestart=true
user=grafana
