# puppet/zulip_ops/templates/supervisor/conf.d/munin_tunnels.conf.erb
; This file was auto-generated by Puppet.  Do not edit by hand.

<%
i = 0
@hosts.each do |host|
-%>
[program:munin-tunnel-<%= host %>]
command=autossh -2 -N -M <%= 20000 + 2 * i %> -L <%= 5000 + i %>:localhost:4949 nagios@<%= host %>.<%= @hosts_domain %>
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=nagios                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/zulip/autossh.<%= host %>.log         ; stdout log path, NONE for none; default AUTO

<%
i += 1
end
-%>
