# puppet/zulip/templates/supervisor/go-camo.conf.erb
[program:go-camo]
command=<%= @bin %> --listen=0.0.0.0:9292 -H "Strict-Transport-Security: max-age=15768000" -H "X-Frame-Options: DENY" --verbose
environment=GOCAMO_HMAC="<%= @camo_key %>"
priority=15
autostart=true
autorestart=true
user=zulip
redirect_stderr=true
stdout_logfile=/var/log/zulip/camo.log
