# servers/puppet/modules/humbug/files/supervisord/conf.d/feedback-bot.conf
[program:feedback-bot]
command=/home/humbug/humbug/api/bots/feedback-bot
autostart=true
user=humbug
