# servers/puppet/modules/humbug/files/supervisord/conf.d/stats.conf
; Supervisor config file for stats.humbughq.com
; It runs graphite, graphiti, statsd


[program:humbug-carbon-cache]
command=python /opt/graphite/bin/carbon-cache.py --debug start
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=root                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/carbon-cache.log         ; stdout log path, NONE for none; default AUTO
directory=/home/humbug/

[program:humbug-carbon-aggregator]
command=python /opt/graphite/bin/carbon-aggregator.py --debug start
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=root                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/carbon-aggregator.log         ; stdout log path, NONE for none; default AUTO
directory=/home/humbug/

[program:humbug-statsd]
command=node stats.js /home/humbug/humbug/servers/puppet/modules/humbug/files/statsd/local.js
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=humbug                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/statsd.log         ; stdout log path, NONE for none; default AUTO
directory=/home/humbug/statsd


[program:humbug-graphiti]
command=/home/humbug/humbug-deployments/current/tools/run-graphiti
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=humbug                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/graphiti.log         ; stdout log path, NONE for none; default AUTO
directory=/home/humbug/graphiti

[group:humbug-stats]
programs=humbug-carbon-cache,humbug-carbon-aggregator,humbug-statsd,humbug-graphiti
