# puppet/zulip_ops/files/supervisor/conf.d/grafana.conf
[program:grafana]
command=/srv/grafana/bin/grafana-server -config /etc/grafana/grafana.ini
directory=/srv/grafana
priority=15
autostart=true
autorestart=true
user=grafana
