# puppet/zulip_ops/files/supervisor/conf.d/postfix.conf
[program:postfix]
command=/usr/sbin/postfix -c /etc/postfix start
