# puppet/zulip_ops/files/supervisor/conf.d/prometheus_node_exporter.conf
[program:prometheus_node_exporter]
command=/usr/local/bin/node_exporter
priority=10
autostart=true
autorestart=true
user=prometheus
redirect_stderr=true
stdout_logfile=/var/log/node_exporter.log
