# puppet/zulip_ops/files/supervisor/conf.d/smokescreen.conf
[program:smokescreen]
command=/usr/local/bin/smokescreen
priority=15
autostart=true
autorestart=true
user=zulip
redirect_stderr=true
stdout_logfile=/var/log/zulip/smokescreen.log
