# puppet/zulip_ops/files/supervisor/conf.d/stats.conf
; Supervisor config file for stats.zulip.net
; It runs graphite, graphiti, statsd


[program:zulip-carbon-cache]
command=/opt/graphite/bin/carbon-cache.py --debug start
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=root                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/zulip/carbon-cache.log         ; stdout log path, NONE for none; default AUTO
directory=/home/zulip/

[program:zulip-carbon-aggregator]
command=/opt/graphite/bin/carbon-aggregator.py --debug start
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=root                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/zulip/carbon-aggregator.log         ; stdout log path, NONE for none; default AUTO
directory=/home/zulip/

[program:zulip-statsd]
command=node stats.js /home/zulip/zulip/puppet/zulip/files/statsd/local.js
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=zulip                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/zulip/statsd.log         ; stdout log path, NONE for none; default AUTO
directory=/home/zulip/statsd


[program:zulip-graphiti]
command=/home/zulip/zulip/tools/run-graphiti
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=zulip                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/zulip/graphiti.log         ; stdout log path, NONE for none; default AUTO
directory=/home/zulip/graphiti

[group:zulip-stats]
programs=zulip-carbon-cache,zulip-carbon-aggregator,zulip-statsd,zulip-graphiti
