# servers/puppet/modules/zulip/files/supervisord/conf.d/zmirror.conf
[program:zmirror]
command=python /home/humbug/humbug/bots/zephyr_mirror.py --root-path=/home/humbug/ --user=tabbott/extra --enable-log=/home/humbug/all_zephyrs_log --forward-class-messages --enable-log=/home/humbug/logs/mirror-log
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=humbug                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/zmirror.log         ; stdout log path, NONE for none; default AUTO
environment=HOME="/home/humbug",USER="humbug"

[program:sync-public-streams]
command=python /home/humbug/humbug/bots/sync-public-streams
priority=200                   ; the relative start priority (default 999)
autostart=true                 ; start at supervisord start (default: true)
autorestart=true               ; whether/when to restart (default: unexpected)
stopsignal=TERM                 ; signal used to kill process (default TERM)
stopwaitsecs=30                ; max num secs to wait b4 SIGKILL (default 10)
user=humbug                    ; setuid to this UNIX account to run the program
redirect_stderr=true           ; redirect proc stderr to stdout (default false)
stdout_logfile=/var/log/humbug/sync-public-streams.log         ; stdout log path, NONE for none; default AUTO
environment=HOME="/home/humbug"
