# puppet/zulip/files/supervisor/conf.d/cron.conf
; Used in dockervoyager.pp, starts cron (for docker use only)
[program:cron]
command = /usr/sbin/cron -f -L 15
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = true
