# puppet/zulip/files/supervisor/conf.d/nginx.conf
; Used in dockervoyager.pp, starts nginx webserver (for docker use only)
[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = true
